# services/pygeoapi_master/demo.pygeoapi.io.cron
0 1 * * * /bin/rm -f /tmp/pygeoapi-process-manager.db > /proc/1/fd/1 2>/proc/1/fd/2
