# services/pygeoapi_master/demo_pygeoapi_io.cron
0 1 * * * /bin/rm -f /tmp/pygeoapi-process-manager.db > /dev/null 2>&1
