# services/update/pygeoapi_update.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin
*/15 * * * *   padmin    /home/padmin/git/services/update/cron.sh > /home/padmin/cron.log 2>&1
