# services/update/pygeoapi_update.cron
*/15 * * * *   padmin    /home/padmin/git/services/update/cron.sh
